# roles/tpm/templates/tpm-updater.service.j2
[Unit]
Description=tpm updater

[Service]
WorkingDirectory={{ tpm_dir }}
ExecStart=/bin/bash -c '/usr/bin/git checkout $(/usr/bin/git describe --tags)'

[Install]
WantedBy=multi-user.target
