# roles/pyenv/files/pyenv-updater.timer
[Unit]
Description=Run pyenv weekly updater

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
