# roles/updater/templates/updater.timer.j2
[Unit]
Description=Run {{ updater_name }} weekly updater

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
