# roles/version_managers/templates/updater.timer.j2
[Unit]
Description=Run {{ item }} weekly updater

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
