# roles/version_manager/templates/updater.timer.j2
[Unit]
Description=Run {{ version_manager_name }} weekly updater

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
